# src/daemon/plasma-knighttimed.service.in
# SPDX-FileCopyrightText: None
# SPDX-License-Identifier: CC0-1.0

[Unit]
Description=KDE Night Time Service
PartOf=graphical-session.target

[Service]
ExecStart=@KDE_INSTALL_FULL_LIBEXECDIR@/knighttimed
BusName=org.kde.NightTime
Slice=background.slice
